# deployments/amazon-ec2/large/scripts/monitoring.conf
description "Samsara Monitor container"
author "Bruno"
start on runlevel [!2345]
stop on runlevel [!2345]
respawn
pre-start exec /usr/bin/docker rm monitoring | true
exec /usr/bin/docker run --name monitoring \
       -p 5555:5555 \
       -p 5556:5556 \
       -p 8083:8083 \
       -p 8086:8086 \
       -p 15000:80 \
       -v /logs/monitoring:/logs \
       -v /data/monitoring:/data \
       `curl "http://169.254.169.254/latest/user-data"` \
       samsara/monitoring

pre-stop script
        /usr/bin/docker stop monitoring
        /usr/bin/docker rm monitoring
end script
